# .docker/config/supervisor/conf.d/crontab.conf
[program:crontab]
command=cron -f
autostart=true
autorestart=true
